# systemd/.config/systemd/user/wsl-wezterm-mux-vsock@.service
[Unit]
Description=WSL VSOCK WezTerm Mux Forwarding %i
ConditionVirtualization=wsl
Requires=wezterm-mux-server.service
After=wezterm-mux-server.service

[Service]
ExecStart=/usr/bin/socat -b65536 VSOCK-LISTEN:%i,fork UNIX-CONNECT:%t/wezterm/sock

[Install]
WantedBy=default.target
